# scripts/configure-leds/configure-leds.service
[Unit]
Description=Configure LEDS

[Service]
Type=oneshot
ExecStart=/bin/sh /root/configure-leds.sh

[Install]
WantedBy=multi-user.target
